# systemd/theos-chariot.service
[Unit]
Description=THEOS Chariot Webhook Server
After=network.target

[Service]
Type=simple
User=tig0_0bitties
WorkingDirectory=/mnt/Covenant/Theos/chariot-repo
Environment=NODE_ENV=production
EnvironmentFile=/mnt/Covenant/Theos/chariot-repo/.env
ExecStart=/usr/bin/node /mnt/Covenant/Theos/chariot-repo/webhook-server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
